# ncid_service
# chkconfig: 3 99 1
# description: CallAlert
# Just put this file into /etc/init.d/
# sudo update-rc.d ncid defaults 98 02

DISPLAY=":0.0"
export DISPLAY

case "$1" in
'start')
xhost + localhost > /dev/null
sudo python /PATH_TO_SCRIPT/ncid-client-py/ncid-client.py start
;;

'stop')
sudo python /PATH_TO_SCRIPT/ncid-client-py/ncid-client.py stop
;;

'status')
sudo python /PATH_TO_SCRIPT/ncid-client-py/ncid-client.py status
;;

'restart')
sudo python /PATH_TO_SCRIPT/ncid-client-py/ncid-client.py restart
;;


esac


